# deploy/ansible/roles/anton-agent/templates/anton-agent.service.j2
[Unit]
Description=anton.computer agent
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={{ anton_user }}
Group={{ anton_user }}
EnvironmentFile={{ anton_install_dir }}/agent.env
WorkingDirectory={{ anton_repo_dir }}
ExecStart=/usr/bin/node packages/agent-server/dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anton-agent

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths={{ anton_install_dir }}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
